# package-root/DEBIAN/prerm
#!/bin/bash
set -e

systemctl stop update-sudoers.service >/dev/null 2>&1 || true
systemctl disable update-sudoers.service >/dev/null 2>&1 || true

exit 0